# infrastructure/targets/colima/k3s/configs/com.lamina.k3s.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lamina.k3s</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/k3s</string>
        <string>server</string>
        <string>--data-dir=/opt/lamina/k3s</string>
        <string>--disable=traefik</string>
        <string>--disable=servicelb</string>
        <string>--disable=metrics-server</string>
        <string>--cluster-cidr=10.42.0.0/16</string>
        <string>--service-cidr=10.43.0.0/16</string>
        <string>--flannel-backend=host-gw</string>
        <string>--kubelet-arg=--max-pods=200</string>
        <string>--kubelet-arg=--kube-reserved=cpu=2000m,memory=8Gi</string>
        <string>--kubelet-arg=--system-reserved=cpu=4000m,memory=32Gi</string>
        <string>--kubelet-arg=--eviction-hard=memory.available&lt;16Gi</string>
        <string>--kubelet-arg=--serialize-image-pulls=false</string>
        <string>--kubelet-arg=--cpu-manager-policy=static</string>
        <string>--kube-apiserver-arg=--max-requests-inflight=800</string>
        <string>--kube-apiserver-arg=--max-mutating-requests-inflight=400</string>
        <string>--etcd-arg=--quota-backend-bytes=10737418240</string>
        <string>--etcd-arg=--auto-compaction-retention=30m</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/opt/lamina/k3s/k3s.log</string>
    
    <key>StandardErrorPath</key>
    <string>/opt/lamina/k3s/k3s.error.log</string>
    
    <key>UserName</key>
    <string>root</string>
    
    <key>GroupName</key>
    <string>wheel</string>
</dict>
</plist>